# FerretDB/setting-up-bacalhau-network-aws-ec2-instance/instance/ferretdb/scripts/ferret-startup.service
[Unit]
Description=FerretDB Startup Script
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/startup.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
